# document-portal/xdg-document-portal.service.in
# SPDX-License-Identifier: LGPL-2.1-or-later
# SPDX-FileCopyrightText: Copyright © the xdg-desktop-portal contributors

[Unit]
Description=flatpak document portal service
PartOf=graphical-session.target
Requires=dbus.service
After=dbus.service

[Service]
BusName=org.freedesktop.portal.Documents
ExecStart=@libexecdir@/xdg-document-portal
Type=dbus
Slice=session.slice
