# document-portal/xdg-permission-store.service.in
# SPDX-License-Identifier: LGPL-2.1-or-later
# SPDX-FileCopyrightText: Copyright © the xdg-desktop-portal contributors

[Unit]
Description=sandboxed app permission store
PartOf=graphical-session.target
Requires=dbus.service
After=dbus.service

[Service]
BusName=org.freedesktop.impl.portal.PermissionStore
ExecStart=@libexecdir@/xdg-permission-store
Type=dbus
Slice=session.slice
