# src/xdg-desktop-portal.service.in
[Unit]
Description=Portal service
PartOf=graphical-session.target
Requisite=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.portal.Desktop
ExecStart=@libexecdir@/xdg-desktop-portal
Slice=session.slice
