# src/xdg-desktop-portal-rewrite-launchers.service.in
# SPDX-License-Identifier: LGPL-2.1-or-later
# SPDX-FileCopyrightText: Copyright © the xdg-desktop-portal contributors

[Unit]
Description=Rewrite dynamic launcher portal entries
# Do it before the graphical session in case the launchers we need to rewrite
# are configured to launch automatically
PartOf=graphical-session-pre.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/xdg-desktop-portal-rewrite-launchers
Slice=session.slice

[Install]
WantedBy=graphical-session-pre.target
